# cron.d/apache-footer
* * * * * root /data/mirror/bin/gen-footer.sh
