# cron.d/update.slave
#####################################################
# Slave mirror crontab
#

MAILTO=mirrors@example.org

* * * * * root /data/mirror/bin/trigger-run

